# deploy/conf/fanmo.cron
0 */3 * * *   ubuntu    cd /home/ubuntu/memberships; bash /home/ubuntu/memberships/prod.deploy.sh 1> /dev/null 2> /home/ubuntu/logs/deploy-cron.err
