# deploy/conf/fanmo.cron
PATH=/usr/bin:/usr/local/bin

0 */1 * * *   ubuntu    cd /home/ubuntu/fanmo && ./prod.deploy.sh 1>> /home/ubuntu/logs/deploy-cron.log 2>> /home/ubuntu/logs/deploy-cron.err
